# rivendel/configuracoes_usando_stow/config/.config/systemd/user/gondor_analise_dados.timer
[Unit]
Description=Backup de diretorio

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=basic.target
