# rivendel/configuracoes_usando_stow/config/.config/systemd/user/gondor_home.timer
[Unit]
Description=Backup do diretório HOME

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=basic.target