# Mordor/config/systemd/user/sincroniza_dropbox_onedrive.timer
[Unit]
Description=Sincronização entre Dropbox e OneDrive

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=basic.target