# rivendel/stow_configs/config/.config/systemd/user/backup_FOTOS.timer
[Unit]
Description=Backup de FOTOS para HD externo

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=basic.target
