# rivendel/stow_configs/config/.config/systemd/user/backup_TARDIS.timer
[Unit]
Description=Backup do HOME para HD externo

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=basic.target
