# rivendel/stow_configs/config/.config/systemd/user/hibernacao.service
[Unit]
Description=Hiberna quando a tela é fechada

[Service]
Type=oneshot
ExecStart=/home/jfreitas/.bin/hibernar.sh

[Install]
WantedBy=multi-user.target
