# Manjaro/config/autostart/key_map.desktop
[Desktop Entry]
Type=Application
Exec=/home/jfreitas/.bin/key_map.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[pt_BR]=key_map
Name=key_map
Comment[pt_BR]=Converte ponto em vírgula no teclado numérico
Comment=Converte ponto em vírgula no teclado numérico
