# Manjaro/config/systemd/user/gondor_bitbucket.timer
[Unit]
Description=Backup dos repositórios do Bitbucket

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=basic.target