# Manjaro/config/systemd/user/gondor_video_aulas.timer
[Unit]
Description=Backup das videos aulas para o Dropbox

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=basic.target