# numenor/dot_config/config/.config/systemd/gondor_root.timer
[Unit]
Description=Agendamento de backup do diretório etc

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=basic.target