# Debian/config/systemd/user/hibernacao.timer
[Unit]
Description=Hiberna quando a tela é fechada

[Timer]
OnCalendar=*:*:0/10
Persistent=true

[Install]
WantedBy=basic.target
