# Gondor/config/systemd/user/teste.service
[Unit]
Description=Testando service

[Service]
Type=oneshot
ExecStart=/bin/bash /home/jfreitas/scripts/teste.sh

[Install]
WantedBy=multi-user.target