# Gondor/config/systemd/user/teste.timer
[Unit]
Description=Testando

[Timer]
OnCalendar=*:0/1
Persisten=true

[Install]
WantedBy=basic.target