# Arch/config/systemd2/arch_keyboard.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/jfreitas/.bin/avell.sh
Type=simple

[Install]
WantedBy=multi-user.target
