# Manjaro-KDE/config/systemd/user/sincroniza_dropbox-onedrive.timer
[Unit]
Description=Sincronização entre Dropbox e OneDrive

[Timer]
OnCalendar=*-*-* 14:20:00
Persistent=true

[Install]
WantedBy=basic.target
